# conf/supervisord.conf
[program:flask_customers]
command = /<ABSOLUTE_PATH>/flask-customers/env/bin/gunicorn customers:app
directory = /<ABSOLUTE_PATH>/flask-customers
user = <USERNAME>
